# scripts/boinc-mobile-helper.service
[Unit]
Description=BOINC Mobile Discovery Helper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.local/share/boinc-mobile/boinc_mobile_helper.py
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
